# src/_raw/usr/local/etc/pdm-wba/cnf/podman-auto-update.timer.d/override.conf
[Timer]
OnCalendar=
RandomizedDelaySec=
OnCalendar=*:0/5
Persistent=true