# src/_raw/etc/systemd/system/pdm-wba-install.path
[Unit]
Description=Monitor for WA install requests

[Path]
PathChanged=/home/wa/install
Unit=pdm-wba-install.service

[Install]
WantedBy=multi-user.target