# src/_raw/etc/systemd/system/pdm-wba-auto-enable.service
# pdm-wba-auto-enable.service
[Unit]
Description=Automatically enables new services generated by Podman
Requires=pdm-wba-service.path
After=pdm-wba-service.path

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 10
ExecStart=/usr/local/bin/pdm-wba-auto-enable-service.sh %f
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target