# src/_raw/etc/systemd/system/pdm-wba-install.path
[Unit]
Description=Monitor fro WA install requests

[Path]
PathChanged=/home/wa/install
PathExists=/home/wa/install

[Install]
WantedBy=multi-user.target