# src/_raw/etc/systemd/system/pdm-wba-service.path
# pdm-wba-service.path
[Unit]
Description=Monitor for WA quadlets change

[Path]
DirectoryNotEmpty=/home/wa/.config/containers/systemd
PathChanged=/home/wa/.config/containers/systemd
PathExists=/home/wa/.config/containers/systemd
Unit=pdm-wba-auto-enable.service

[Install]
WantedBy=multi-user.target