# src/_raw/etc/profile.d/wa.sh
#!/bin/bash
# Sources
source /usr/local/etc/pdm-wba/src/wa_env_nw.sh
source /usr/local/etc/pdm-wba/src/wa_env_aliases.sh
source /usr/local/etc/pdm-wba/src/wa_env.sh
# Publish
wa_env
wa_env_nw
wa_env_aliases